# deployment/conduit.service
# deployment/conduit.service
# --------------------------
# Production systemd unit for Psiphon Conduit 2.0.0 (Issue #45).
#
# Installed to /etc/systemd/system/conduit.service by install.sh (Phase 2x-e).
# Managed by: install.sh, update.sh, uninstall.sh
#
# Binary:    /opt/conduit/conduit   (root:root 755)
# Data dir:  /var/lib/conduit/      (conduit:conduit 700)
# Key file:  /var/lib/conduit/data/conduit_key.json  (mode 0600, auto-created)
#
# Metrics endpoint: http://127.0.0.1:9090/metrics  (loopback only)
#   --metrics-addr MUST be 127.0.0.1:9090 (not :9090) to avoid exposing
#   the Prometheus endpoint on the public interface.
#   The port (9090) must match conduit.metrics_port in config.json.
#
# Defaults:
#   --max-common-clients 50    Adjust for available bandwidth and CPU
#   --bandwidth 40             Upload capacity in Mbps
#   reduced window OFF         CCC_REDUCED_* default to empty/0; when CCC sets a
#                              window, tunnel-core applies the daily HH:MM-UTC
#                              reduced limits at runtime, with no restart at the
#                              window boundaries (BS1).
#
# UFW: Conduit binds additional UDP ports for inproxy traffic.  The ports are
# not fixed in advance.  After install, discover them with:
#   ss -ulnp | grep conduit
# Then add UFW rules for any UDP ports shown.
#
# conduit_key.json: The file at /var/lib/conduit/data/conduit_key.json is the
# node identity keypair.  Losing it resets your Psiphon broker reputation to
# zero.  Back it up before any migration or hardware replacement.

[Unit]
Description=Psiphon Conduit inproxy node
Documentation=https://github.com/Psiphon-Inc/conduit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=conduit
Group=conduit
WorkingDirectory=/var/lib/conduit
# CCC-managed configuration knobs (M2). Defaults live here; CCC overrides them by
# writing ONLY Environment= lines to the drop-in
# /etc/systemd/system/conduit.service.d/ccc.conf (via the root helper
# ccc-apply-conduit-config). ExecStart references the BRACED form ${VAR} so each
# value is substituted as exactly one argument, with no shell. Do not change the
# braced form to $VAR (that would whitespace-split the value).
Environment=CCC_MAX_COMMON_CLIENTS=50
Environment=CCC_BANDWIDTH_MBPS=40
# Reduced-window knobs (BS1). Empty times + zero limits = reduced mode OFF.
# CCC overrides these the same way as the knobs above: ONLY Environment= lines
# in the drop-in, written by the root helper ccc-apply-conduit-config. The
# reduced times are HH:MM (24-hour, UTC); the byte/sec limits are integers.
Environment=CCC_REDUCED_START=
Environment=CCC_REDUCED_END=
Environment=CCC_REDUCED_MAXCOMMON=0
Environment=CCC_REDUCED_UP=0
Environment=CCC_REDUCED_DOWN=0
ExecStart=/opt/conduit/conduit start \
    --data-dir /var/lib/conduit/data \
    --metrics-addr 127.0.0.1:9090 \
    --max-common-clients ${CCC_MAX_COMMON_CLIENTS} \
    --bandwidth ${CCC_BANDWIDTH_MBPS} \
    --set InproxyReducedStartTime=${CCC_REDUCED_START} \
    --set InproxyReducedEndTime=${CCC_REDUCED_END} \
    --set InproxyReducedMaxCommonClients=${CCC_REDUCED_MAXCOMMON} \
    --set InproxyReducedLimitUpstreamBytesPerSecond=${CCC_REDUCED_UP} \
    --set InproxyReducedLimitDownstreamBytesPerSecond=${CCC_REDUCED_DOWN}
Restart=on-failure
RestartSec=10s
ProtectSystem=strict
ReadWritePaths=/var/lib/conduit
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
